# build/app/supervisord/conf.d/WorkerDispatchJob.conf
[program:WorkerDispatchJob]
command=bin/console vtv:worker:start
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/%(program_name)s.log
numprocs=10
startretries=30
directory=/var/www
autostart=true
autorestart=true
user = root